# deploy/supervisor.conf
; =============================================================================
; Supervisor Configuration for SocialDL Queue Workers
; =============================================================================
; Location: /etc/supervisor/conf.d/socialdl.conf
;
; Install: sudo apt install supervisor
; Reload:  sudo supervisorctl reread && sudo supervisorctl update
; Status:  sudo supervisorctl status
; =============================================================================

; ── Default Queue Worker (general tasks) ─────────────────────────────────────
[program:socialdl-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/socialdl/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --memory=128
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/socialdl-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

; ── Download Queue Worker (aria2c downloads) ─────────────────────────────────
[program:socialdl-downloads]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/socialdl/artisan queue:work redis --queue=downloads --sleep=3 --tries=3 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/socialdl-downloads.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

; ── Merge Queue Worker (FFmpeg operations) ───────────────────────────────────
[program:socialdl-merge]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/socialdl/artisan queue:work redis --queue=merge --sleep=3 --tries=2 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/socialdl-merge.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

; ── HLS Queue Worker (HLS generation) ───────────────────────────────────────
[program:socialdl-hls]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/socialdl/artisan queue:work redis --queue=hls --sleep=5 --tries=2 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/socialdl-hls.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

; ── Cleanup Queue Worker ─────────────────────────────────────────────────────
[program:socialdl-cleanup]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/socialdl/artisan queue:work redis --queue=cleanup --sleep=10 --tries=1 --max-time=3600 --memory=64
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/socialdl-cleanup.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stopwaitsecs=120

; ── Laravel Scheduler (cron replacement) ─────────────────────────────────────
[program:socialdl-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do php /var/www/socialdl/artisan schedule:run --no-interaction >> /dev/null 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/socialdl-scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2

; ── Worker Group ─────────────────────────────────────────────────────────────
[group:socialdl-workers]
programs=socialdl-default,socialdl-downloads,socialdl-merge,socialdl-hls,socialdl-cleanup,socialdl-scheduler
priority=999
